# 99-cd-ripper.rules
SUBSYSTEM=="block",KERNEL=="sr[0-9]*",ACTION=="change",RUN+="/bin/systemctl start cd-ripper"
