# cd-ripper.service
[Unit]
Description=Auto CD ripper

[Service]
Type=oneshot
ExecStart=/bin/cd-rip

[Install]
WantedBy=multi-user.target